# files/rc.local
#!/bin/sh -e

if [ "$(grep /var/lib/docker /proc/mounts | grep rootfs)" ]
then
  export DOCKER_RAMDISK=true
fi

udhcpc

exit 0

